# crontab
# m h  dom mon dow   command
0,15,30,45 * * * * ~/Git/dotfiles/commit_files.sh
